# powerline-daemon.service
[Service]
Type=simple
ExecStart=/usr/bin/zsh -c "powerline-daemon --foreground"

[Install]
WantedBy=default.target
