# powerline-daemon.service
[Service]
Type=simple
ExecStart=/usr/bin/zsh -c "$HOME/.local/bin/powerline-daemon --foreground"

[Install]
WantedBy=default.target
